# backfile/.config/autostart/mictray.desktop
[Desktop Entry]
Exec=/usr/bin/mictray
Name=mictray
Type=Application
